# nomorenotes.service
[Unit]
Description=Nomorenotes chat and file sharing
After=network.target

[Service]
User=nathan
Environment=PORT=8954
Environment=URL=https://n.pi.pool.bad.mn
Environment=AUTH0_SECRET=29c4acaf31cdb541e0aae0ca23cffd4e7711946c8917c2f135b8a3b933987e98ec8393b62e92e55fff8bd5b853ea899413cb80291daae134da09ffab259bd31a09a57effe30c104225ea7caa73f05d732b891ad5482b7e5436c1d9797723cddb0df09cc5c5fbb790e570871c329f27ca68e222f51111b0e9b5e7c9cf548b7f6432150665a5f60da1f5b733be882ed8b2c33e2335bcaee037022646786dddbf4c085734c65eef643a7c96f82e1da1e96f68e3530ffa5ca01a5b8c38098e4f62435451edaae239a29ec803b1ad7fca1868c0854ec26c28c16f62fab66994f5f3dc7fd494b26d7f16481336cc39be8f4924278346fad95dae652c3be6aaa966c00e16f7d182e11f8967c3842250ae311112c52a6e1b44ee3c79eebc43c485ba3284f8b05905e009b8178cf24cc450a9c6c2746b139a181ac5fe1bf13d862789f89abd83b2b22e78d36055622bfdadbb89c31e793a0785ea62f34ba0860918ab73565389bfdc86e9283aa161bc4959bce23f71c9ca823da111a3b2673dcc043b56187a91918954e28cc058537406d1bd5c3fc1dcf73513aca696cf952104620d11791a0178f893c003bf2d971ee1ee44560b5321fb28cf3256166cd2a80ad5d878f4de69c7e3079b2ce2318fa472b77e785911a1b7638b685b75f0409f40c3dc4c2bee4029f20e0e4e21348fda0985c13d6c457d9824b4c9f2cb3d831ff87dc5cb44`
Environment=MAIL_URL=https://discord.com/api/webhooks/1050210534034788425/6B0nmrszfeLhZ_OwdizqMAsloWXTC4HKLoi5GTHyZ7wyURPuxKBC8i94RJ9zQvwfKgjV
WorkingDirectory=/home/nathan/nomorenotes/
ExecStart=/usr/local/nvm/versions/node/v19.2.0/bin/node server.js
Restart=always

[Install]
WantedBy=multi-user.target
